# shadowsocks.conf
[program:shadowsocks]
user=root
directory=/usr/local/bin
command=ssserver -p 8379 -k bily -m aes-256-cfb
process_name=%(program_name)s
autostart=true
